# conf/etc/services.d/fotd/run
#!/bin/sh
NAME="FOTD"
COMMAND="python3 /fotd/fotd.py"

exec service-handler "$NAME" "$COMMAND"
